# live/crontab
* * * * * /home/marcel/get_live_viewer.py
0 * * * * /home/marcel/update_live_data.sh
