# live/crontab
* * * * * /home/marcel/rbtv/rbtvstats/live/get_live_data.py
0 * * * * /home/marcel/rbtv/rbtvstats/live/update_live_data.py
